# streamlined.service
[Unit]
Description=Streamlined Dashboard, Frontend & Simulation Suite
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/streamlined
ExecStart=/bin/bash /streamlined/start.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target